# cron.daily/eebydeeby
#!/usr/bin/env bash

exec env -i eebydeeby "$@"
