# systemd/eebydeeby.timer
[Unit]
Description=Timer unit for Eebydeeby.

[Timer]
OnCalendar=Sat *-*-* 04:00:00
WakeSystem=true
