# systemd/eebydeeby.service
[Unit]
Description=Eeebdeeby, the janky automated Gentoo custodian. Beep!

[Service]
Type=exec
ExecStart=eebydeeby
